# jcsm-interne/deploy/jcsm-interne.service
[Unit]
Description=JCSM Interne API
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/jcsm-interne/backend
Environment="PATH=/var/www/jcsm-interne/venv/bin"
EnvironmentFile=/var/www/jcsm-interne/.env
ExecStart=/var/www/jcsm-interne/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
